# netwatch.service
[Unit]
Description=NetWatch - Network Performance Monitor
After=network.target

[Service]
Type=simple
User=YOUR_USERNAME
WorkingDirectory=/path/to/netwatch
Environment="PATH=/path/to/netwatch/.venv/bin"
ExecStart=/path/to/netwatch/.venv/bin/python main.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
